# docker/supervisord-rabbit.conf
; 与 /etc/supervisor/supervisord.conf 配合：官方 rabbitmq 镜像入口脚本 + Spring Boot

[program:rabbitmq]
command=/usr/local/bin/docker-entrypoint.sh rabbitmq-server
autostart=true
autorestart=true
priority=100
startsecs=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:arb-app]
command=/app/wait-rabbit-and-run.sh
user=arb
directory=/app
autostart=true
autorestart=true
priority=200
startretries=5
startsecs=25
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
